# lime/supervisor/lime.conf
[program:lime]
command = /lime-code/lime/start.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/lime.err.log
stdout_logfile=/var/log/lime.out.log